# files/ydotool-user.service
[Unit]
# home-os: Fedora liefert ydotool nur als System-Unit. Voxtype (User-Service)
# erwartet aber ein User-`ydotool.service` (Wants=). Das stellen wir hier bereit:
# ydotoold läuft als User und greift via input-Gruppe auf /dev/uinput zu.
Description=ydotoold (user) – uinput-Tippen/Strg+V für Voxtype paste-Modus
PartOf=graphical-session.target

[Service]
Type=simple
# Socket im User-Runtime-Dir ($XDG_RUNTIME_DIR = %t), nur für den User lesbar.
ExecStart=/usr/bin/ydotoold --socket-path=%t/.ydotool_socket --socket-perm=0600
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
